# scripts/bash/daemon/port_daemon@.service
# /etc/systemd/system/port-daemon@.service
[Unit]
Description=Port daemon (TCP) on port %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=portdaemon
Group=portdaemon
ExecStart=/usr/local/bin/port-daemon.sh %i 0.0.0.0
Restart=always
RestartSec=2

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
